# etc/systemd/system/httpd.service.d/override.conf
[Unit]
After=syslog.target network.target consul-template.service

[Service]
EnvironmentFile=-/etc/sysconfig/httpd.d/*
Restart=on-failure
RestartSec=5s
