# etc/systemd/system/tomcat.service.d/override.conf
[Unit]
Requires=consul-template.service memcached.service
After=syslog.target network.target consul-template.service memcached.service

[Service]
Restart=on-failure
RestartSec=60s
